# pkg/scripts/postinstall
#!/bin/bash
# Claw Web — postinstall (최소화 버전 v1.2.43)
# - 파일 권한 고정 + 데이터 디렉토리 + LaunchAgent 만 수행
# - Claude CLI 설치/인증/터널 설정은 웹 UI Settings 에서 사용자가 직접 처리
set -e

export LANG=en_US.UTF-8
export LC_ALL=en_US.UTF-8
export LANGUAGE=en_US.UTF-8

INSTALL_DIR="/usr/local/lib/claw-web"
LOG="/tmp/claw-web-install.log"
echo "[postinstall] $(date)" >> "$LOG"

export PATH="/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:$PATH"

# ── 사용자 / 홈 감지 ─────────────────────────────────────
REAL_USER=$(stat -f '%Su' /dev/console 2>/dev/null || echo "$SUDO_USER")
REAL_HOME=$(eval echo "~$REAL_USER" 2>/dev/null)
[ -z "$REAL_HOME" ] || [ "$REAL_HOME" = "~" ] && REAL_HOME="$HOME"

# ── Node.js 탐색 (brew / nvm / volta / fnm / asdf) ───────
find_node() {
  for p in \
    /opt/homebrew/bin/node \
    /opt/homebrew/opt/node@22/bin/node \
    /opt/homebrew/opt/node@20/bin/node \
    /opt/homebrew/opt/node/bin/node \
    /usr/local/bin/node \
    /usr/local/opt/node@22/bin/node \
    /usr/local/opt/node@20/bin/node; do
    [ -x "$p" ] && echo "$p" && return 0
  done
  if [ -d "$REAL_HOME/.nvm/versions/node" ]; then
    local latest
    latest=$(ls -1 "$REAL_HOME/.nvm/versions/node" 2>/dev/null | sort -V | tail -1)
    [ -n "$latest" ] && [ -x "$REAL_HOME/.nvm/versions/node/$latest/bin/node" ] && \
      echo "$REAL_HOME/.nvm/versions/node/$latest/bin/node" && return 0
  fi
  [ -x "$REAL_HOME/.volta/bin/node" ] && echo "$REAL_HOME/.volta/bin/node" && return 0
  if [ -d "$REAL_HOME/.local/share/fnm/node-versions" ]; then
    local latest
    latest=$(ls -1 "$REAL_HOME/.local/share/fnm/node-versions" 2>/dev/null | sort -V | tail -1)
    [ -n "$latest" ] && [ -x "$REAL_HOME/.local/share/fnm/node-versions/$latest/installation/bin/node" ] && \
      echo "$REAL_HOME/.local/share/fnm/node-versions/$latest/installation/bin/node" && return 0
  fi
  if [ -d "$REAL_HOME/.asdf/installs/nodejs" ]; then
    local latest
    latest=$(ls -1 "$REAL_HOME/.asdf/installs/nodejs" 2>/dev/null | sort -V | tail -1)
    [ -n "$latest" ] && [ -x "$REAL_HOME/.asdf/installs/nodejs/$latest/bin/node" ] && \
      echo "$REAL_HOME/.asdf/installs/nodejs/$latest/bin/node" && return 0
  fi
  command -v node 2>/dev/null
}

NODE_BIN=$(find_node)
if [ -n "$NODE_BIN" ] && ! "$NODE_BIN" --version >/dev/null 2>&1; then
  NODE_BIN=""
fi
if [ -z "$NODE_BIN" ]; then
  osascript -e 'display dialog "Node.js 를 찾을 수 없습니다.\nhttps://nodejs.org 에서 설치 후 다시 시도하세요." buttons {"확인"} default button 1 with icon stop with title "Claw Web"' 2>/dev/null || true
  echo "[postinstall] FATAL: Node.js not found" >> "$LOG"
  exit 1
fi

NODE_DIR=$(dirname "$NODE_BIN")
USER_DATA="$REAL_HOME/.claw-web"
PLIST_DIR="$REAL_HOME/Library/LaunchAgents"
mkdir -p "$USER_DATA" "$PLIST_DIR"
echo "[postinstall] node=$NODE_BIN home=$REAL_HOME user=$REAL_USER" >> "$LOG"

notify() {
  osascript -e "display notification \"$1\" with title \"Claw Web 설치 중\"" 2>/dev/null || true
  echo "[postinstall] $1" >> "$LOG"
}
as_user() { su "$REAL_USER" -c "$1" 2>> "$LOG" || true; }

# ═══════════════════════════════════════════════════════
# Step 1. 파일 소유권 + 권한 (v1.2.43 핵심 fix)
# ═══════════════════════════════════════════════════════
# pkg payload 가 root:wheel 로 설치되면 user 접근 불가. user 소유로 넘기고
# 모든 파일에 read+execute 권한 보장.
notify "① 파일 준비 중..."
chown -R "$REAL_USER:staff" "$INSTALL_DIR" 2>/dev/null || true
chmod -R a+rX "$INSTALL_DIR" 2>/dev/null || true
echo "[postinstall] permissions set: owner=$REAL_USER mode=a+rX" >> "$LOG"

# 설치 검증 — server/index.js 가 실제 존재하는지
if [ ! -f "$INSTALL_DIR/server/index.js" ]; then
  osascript -e "display dialog \"설치 파일이 손상되었습니다.\n$INSTALL_DIR/server/index.js 를 찾을 수 없습니다.\n다시 설치해주세요.\" buttons {\"확인\"} with icon stop with title \"Claw Web\"" 2>/dev/null || true
  echo "[postinstall] FATAL: server/index.js missing at $INSTALL_DIR" >> "$LOG"
  exit 1
fi

# ═══════════════════════════════════════════════════════
# Step 2. 데이터 디렉토리 + 심볼릭 링크 + Bearer 토큰
# ═══════════════════════════════════════════════════════
# 서버는 REPO_ROOT = INSTALL_DIR 에서 데이터를 읽/쓰기 때문에
# 모든 유저 데이터 파일/디렉토리는 USER_DATA → INSTALL_DIR 로 심링크 걸어야 한다.
# (누락된 항목은 업그레이드 때 rm -rf 로 사라짐)
for f in sessions.json projects.json skills.json web-metadata.json secrets.json \
         hooks.json schedules.json accounts.json push-subscriptions.json; do
  if [ ! -f "$USER_DATA/$f" ]; then
    case "$f" in
      secrets.json) echo '{"version":1,"backends":{}}' > "$USER_DATA/$f"; chmod 600 "$USER_DATA/$f" ;;
      sessions.json) echo '{"version":1,"sessions":{}}' > "$USER_DATA/$f" ;;
      projects.json) echo '{"version":1,"projects":[]}' > "$USER_DATA/$f" ;;
      skills.json) echo '{"version":1,"skills":{}}' > "$USER_DATA/$f" ;;
      web-metadata.json) echo '{"version":1,"agents":{}}' > "$USER_DATA/$f" ;;
      accounts.json) echo '{"version":1,"accounts":{}}' > "$USER_DATA/$f" ;;
      push-subscriptions.json) echo '{"version":1,"subscriptions":[]}' > "$USER_DATA/$f" ;;
      *) echo '{"version":1}' > "$USER_DATA/$f" ;;
    esac
  fi
  # 기존 payload 에 남아있을 수 있는 실 파일을 지우고 심링크로 교체
  [ -e "$INSTALL_DIR/$f" ] && [ ! -L "$INSTALL_DIR/$f" ] && rm -f "$INSTALL_DIR/$f"
  ln -sf "$USER_DATA/$f" "$INSTALL_DIR/$f" 2>/dev/null || true
done

# 기존 web-config.json 에 토큰이 있으면 그 값을 유지, 없을 때만 1234 사용
EXISTING_TOKEN=""
if [ -f "$USER_DATA/web-config.json" ] && [ -n "$NODE_BIN" ]; then
  EXISTING_TOKEN=$("$NODE_BIN" -e "
    try {
      const d = JSON.parse(require('fs').readFileSync(process.argv[1], 'utf8'));
      const t = d && d.auth && d.auth.token;
      if (typeof t === 'string' && t.length > 0) process.stdout.write(t);
    } catch (e) {}
  " "$USER_DATA/web-config.json" 2>/dev/null || true)
fi
BEARER_TOKEN="${EXISTING_TOKEN:-1234}"
# agents-config.json 이 배포 payload 에 없을 수도 있으므로 기본 구조로 생성
for f in backends.json agents-config.json; do
  if [ ! -f "$USER_DATA/$f" ]; then
    if [ -f "$INSTALL_DIR/$f" ] && [ ! -L "$INSTALL_DIR/$f" ]; then
      cp "$INSTALL_DIR/$f" "$USER_DATA/$f"
    else
      case "$f" in
        agents-config.json) echo '{"agents":{},"channels":{}}' > "$USER_DATA/$f" ;;
        backends.json) echo '{"version":1,"backends":{}}' > "$USER_DATA/$f" ;;
      esac
    fi
  fi
  [ -e "$INSTALL_DIR/$f" ] && [ ! -L "$INSTALL_DIR/$f" ] && rm -f "$INSTALL_DIR/$f"
  ln -sf "$USER_DATA/$f" "$INSTALL_DIR/$f" 2>/dev/null || true
done

# 데이터 디렉토리 — 세션 본문, 업로드, 로그, 백업
# 서버가 쓰는 경로 자체를 심링크로 돌려서 모든 write 가 USER_DATA 아래로 가게 한다.
for d in sessions-store agents-config-backups logs uploads backups; do
  mkdir -p "$USER_DATA/$d" 2>/dev/null || true
  # INSTALL_DIR 에 실 디렉토리가 있으면 제거 (payload 에서 같이 왔을 수 있음)
  if [ -d "$INSTALL_DIR/$d" ] && [ ! -L "$INSTALL_DIR/$d" ]; then
    rm -rf "$INSTALL_DIR/$d"
  fi
  ln -sfn "$USER_DATA/$d" "$INSTALL_DIR/$d" 2>/dev/null || true
done

if [ ! -f "$USER_DATA/web-config.json" ]; then
  cat > "$USER_DATA/web-config.json" <<CFG
{
  "port": 3838,
  "allowedRoots": [],
  "features": {"dashboard":true,"agentsPage":true,"dragAndDrop":true,"chat":true,"activityFeed":true,"mdEditor":true},
  "auth": {"enabled": true, "token": "$BEARER_TOKEN"},
  "configPath": "$USER_DATA/agents-config.json"
}
CFG
else
  # 기존 설정 보존 — 토큰/auth.enabled 는 사용자가 변경한 값을 그대로 유지하고
  # 비어있을 때만 기본값을 채운다. configPath 가 누락되면 보정.
  "$NODE_BIN" -e "
    const fs = require('fs');
    const path = process.argv[1], token = process.argv[2], userData = process.argv[3];
    const d = JSON.parse(fs.readFileSync(path, 'utf8'));
    if (!d.auth) d.auth = {};
    if (typeof d.auth.token !== 'string' || d.auth.token.length === 0) {
      d.auth.token = token;
    }
    if (typeof d.auth.enabled !== 'boolean') {
      d.auth.enabled = true;
    }
    if (!d.configPath || typeof d.configPath !== 'string') {
      d.configPath = userData + '/agents-config.json';
    }
    fs.writeFileSync(path, JSON.stringify(d, null, 2));
  " "$USER_DATA/web-config.json" "$BEARER_TOKEN" "$USER_DATA" 2>/dev/null || true
fi
ln -sf "$USER_DATA/web-config.json" "$INSTALL_DIR/web-config.json" 2>/dev/null || true
chown -R "$REAL_USER" "$USER_DATA" 2>/dev/null || true
echo "[postinstall] data dir ready, token=$BEARER_TOKEN" >> "$LOG"

# ═══════════════════════════════════════════════════════
# Step 3. 작업 디렉토리 선택 (유일하게 남은 대화상자)
# ═══════════════════════════════════════════════════════
WORK_DIR=$(osascript 2>/dev/null <<OSA
try
  set chosen to choose folder with prompt "에이전트가 파일을 읽고 쓸 프로젝트 폴더를 선택하세요." default location (POSIX file "$REAL_HOME")
  return POSIX path of chosen
on error
  return ""
end try
OSA
)
WORK_DIR="${WORK_DIR:-$REAL_HOME}"
WORK_DIR=$(echo "$WORK_DIR" | sed "s|^~|$REAL_HOME|")
mkdir -p "$WORK_DIR" 2>/dev/null || WORK_DIR="$REAL_HOME"

"$NODE_BIN" -e "
  const fs = require('fs');
  const [p, w, i] = process.argv.slice(1);
  const d = JSON.parse(fs.readFileSync(p, 'utf8'));
  const roots = d.allowedRoots || [];
  for (const r of [w, i]) {
    const rs = r.replace(/\/*$/, '/');
    if (!roots.some(e => e.replace(/\/*$/,'/') === rs)) roots.push(rs);
  }
  d.allowedRoots = roots;
  fs.writeFileSync(p, JSON.stringify(d, null, 2));
" "$USER_DATA/web-config.json" "$WORK_DIR" "$INSTALL_DIR" 2>/dev/null || true
echo "[postinstall] work_dir=$WORK_DIR" >> "$LOG"

# ═══════════════════════════════════════════════════════
# Step 4. LaunchAgent 등록 + 서버 기동
# ═══════════════════════════════════════════════════════
notify "② 자동 시작 등록 중..."
cat > "$PLIST_DIR/com.claw-web.server.plist" <<EOFPLIST
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.claw-web.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>$NODE_BIN</string>
        <string>$INSTALL_DIR/server/index.js</string>
    </array>
    <key>WorkingDirectory</key><string>$INSTALL_DIR</string>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>StandardOutPath</key><string>/tmp/claw-web.log</string>
    <key>StandardErrorPath</key><string>/tmp/claw-web.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key><string>$NODE_DIR:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key><string>$REAL_HOME</string>
        <key>NODE_ENV</key><string>production</string>
    </dict>
</dict>
</plist>
EOFPLIST

chown "$REAL_USER" "$PLIST_DIR/com.claw-web.server.plist" 2>/dev/null || true

REAL_UID=$(id -u "$REAL_USER" 2>/dev/null || echo "")
launchctl bootout "gui/$REAL_UID/com.claw-web.server" 2>/dev/null || \
  as_user "launchctl unload '$PLIST_DIR/com.claw-web.server.plist'" 2>/dev/null || true
if [ -n "$REAL_UID" ]; then
  launchctl bootstrap "gui/$REAL_UID" "$PLIST_DIR/com.claw-web.server.plist" 2>/dev/null || \
    as_user "launchctl load '$PLIST_DIR/com.claw-web.server.plist'" 2>/dev/null || true
  launchctl kickstart -k "gui/$REAL_UID/com.claw-web.server" 2>/dev/null || true
else
  as_user "launchctl load '$PLIST_DIR/com.claw-web.server.plist'" 2>/dev/null || true
fi

# fallback: LaunchAgent 실패 시 직접 기동
SERVER_UP=0
for i in $(seq 1 8); do
  sleep 1
  if lsof -i :3838 -sTCP:LISTEN >/dev/null 2>&1; then SERVER_UP=1; break; fi
done
if [ "$SERVER_UP" = "0" ]; then
  echo "[postinstall] LaunchAgent failed — fallback direct start" >> "$LOG"
  su - "$REAL_USER" -c "export NODE_ENV=production; export HOME='$REAL_HOME'; cd '$INSTALL_DIR'; nohup '$NODE_BIN' '$INSTALL_DIR/server/index.js' >> /tmp/claw-web.log 2>&1 & disown" 2>> "$LOG" || true
fi

# ═══════════════════════════════════════════════════════
# Step 5. 완료 다이얼로그 + 브라우저 오픈
# ═══════════════════════════════════════════════════════
echo "[postinstall] waiting for port 3838..." >> "$LOG"
for i in $(seq 1 15); do
  sleep 1
  lsof -i :3838 -sTCP:LISTEN >/dev/null 2>&1 && { echo "[postinstall] port 3838 open after ${i}s" >> "$LOG"; break; }
done

SUMMARY="Claw Web 설치 완료!

🔑 Bearer 토큰 (외부 접속용):
$BEARER_TOKEN

📁 작업 디렉토리:
$WORK_DIR

🏠 접속: http://localhost:3838

Claude CLI 및 원격 접속 설정은 웹에서 Settings 메뉴를 통해 처리하세요."

echo "$BEARER_TOKEN" > "$USER_DATA/.bearer-token"
chown "$REAL_USER" "$USER_DATA/.bearer-token" 2>/dev/null || true

osascript 2>/dev/null -e "display dialog \"$SUMMARY\" buttons {\"토큰 복사\",\"브라우저 열기\"} default button 2 with icon note with title \"Claw Web 설치 완료\"" > /tmp/claw-web-finish.txt || true
RESULT=$(cat /tmp/claw-web-finish.txt 2>/dev/null)
rm -f /tmp/claw-web-finish.txt

if echo "$RESULT" | grep -q "토큰 복사"; then
  as_user "echo -n '$BEARER_TOKEN' | pbcopy"
fi
as_user "open http://localhost:3838"

echo "[postinstall] COMPLETE (token=$BEARER_TOKEN)" >> "$LOG"
exit 0
